# home.linux/.config/systemd/user/xremap.service
[Unit]
Description=Xremap
After=default.target

[Service]
ExecStart=/usr/bin/xremap --watch=device --ignore=solaar --ignore=ydotoold --ignore=Espanso %h/.config/xremap/config.yml
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
