# roles/pluto/files/etc/systemd/system/tuxedo-charging-profile.timer
[Unit]
Description=Configure Tuxedo charging profile automatically

[Timer]
OnBootSec=30
OnUnitActiveSec=600

[Install]
WantedBy=timers.target
