# roles/pluto/files/etc/systemd/system/tuxedo-charging-profile.service
[Unit]
Description=Configure charging profile

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tuxedo-charging-profile
