# roles/pluto/files/etc/udev/rules.d/60-tuxedo-kbd-leds-uaccess.rules
DRIVERS=="ite_8291", SUBSYSTEM=="leds", RUN+="/usr/bin/chmod a+w %S%p/brightness %S%p/multi_intensity"
