# roles/pluto/files/etc/udev/rules.d/99-pluto-leds-permissions.rules
KERNEL=="rgb:kbd_backlight", SUBSYSTEMS=="leds", RUN+="/usr/bin/chmod a+w /sys/class/leds/rgb:kbd_backlight/brightness"
